# archiso/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
#ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin avery - $TERM
ExecStart=-/sbin/agetty --autologin avery --noclear %I 38400 linux
